# supervisord-all-in-one.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# MongoDB Service
[program:mongodb]
command=mongod --dbpath /data/db --bind_ip_all --port 27017
user=mongodb
autorestart=true
stdout_logfile=/var/log/mongodb/mongodb.log
stderr_logfile=/var/log/mongodb/mongodb_error.log
priority=100
startsecs=10
stopwaitsecs=30

# Backend API Service
[program:backend]
command=/app/backend/venv/bin/python -m uvicorn enhanced_server:app --host 0.0.0.0 --port 8001
directory=/app/backend
user=root
environment=PATH="/app/backend/venv/bin",PYTHONPATH="/app/backend"
autorestart=true
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_error.log
priority=200
startsecs=15
stopwaitsecs=10
depends_on=mongodb

# Nginx Service
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_error.log
priority=300
startsecs=5
stopwaitsecs=10
depends_on=backend

# Log rotation
[program:logrotate]
command=/usr/sbin/logrotate -f /etc/logrotate.d/portfolio
user=root
autorestart=false
startsecs=0
exitcodes=0

[group:portfolio]
programs=mongodb,backend,nginx
priority=999